# infra/systemd/family-messenger-backend.service
[Unit]
Description=Family Messenger Backend
After=network.target

[Service]
Type=simple
User=family
Group=family
WorkingDirectory=/opt/family-messenger
EnvironmentFile=/etc/family-messenger/backend.env
ExecStart=/usr/bin/java -Dktor.deployment.port=${BACKEND_PORT} -jar /opt/family-messenger/family-messenger-backend-all.jar
Restart=always
RestartSec=5
SuccessExitStatus=143

[Install]
WantedBy=multi-user.target
